# templates/packer/bitcoind/bitcoind.service.erb
[Unit]
Description=The Bitcoin Blockchain Command Line Interface
After=network.target

[Service]
User=bitcoin
Group=bitcoin
ExecStart=/usr/local/sbin/bitcoind
[Install]
WantedBy=default.target
